# infrastructure/rxf-admin.service
# Optional: bare-metal systemd unit for running the backend without Docker.
# Place at /etc/systemd/system/rxf-admin.service inside the LXC if you do not
# want to use docker-compose.  In a docker-compose deployment you do NOT need
# this file — Docker's restart policy already handles supervision.

[Unit]
Description=rxf-sys admin dashboard backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rxfadmin
Group=rxfadmin
WorkingDirectory=/opt/rxf-admin/backend
EnvironmentFile=/opt/rxf-admin/infrastructure/.env
ExecStart=/opt/rxf-admin/backend/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8080 \
    --proxy-headers \
    --forwarded-allow-ips=*
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rxf-admin
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
